# deploy/lingo-dungeon.service
[Unit]
Description=Lingo Dungeon – Flask app (temporary HTTP, no prod WSGI)
After=network.target

[Service]
Type=simple
User=claude-agent
WorkingDirectory=/home/claude-agent/lingo_pkmn
ExecStart=/home/claude-agent/lingo_pkmn/deploy/run-flask.sh
Restart=always
RestartSec=5

Environment="FLASK_RUN_HOST=127.0.0.1"
Environment="FLASK_RUN_PORT=5000"
Environment="FLASK_DEBUG=0"
# Generate a token: python3 -c "import secrets; print(secrets.token_hex(32))"
Environment="LEADERBOARD_ADMIN_TOKEN=a4a6e7dab0648f074a61757218954cc96e8767f1db25e58bee135f86f13eb919"

[Install]
WantedBy=multi-user.target
